# deploy/linux/systemd/vox-worker.service
[Unit]
Description=VOX background worker
After=network-online.target vox-web.service
Wants=network-online.target

[Service]
Type=simple
User=vox
Group=vox
WorkingDirectory=/opt/vox
EnvironmentFile=/opt/vox/.env
ExecStart=/opt/vox/.venv/bin/python /opt/vox/worker.py
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
