# deploy/linux/systemd/vox-maintenance.service
[Unit]
Description=VOX maintenance cleanup

[Service]
Type=oneshot
User=vox
Group=vox
WorkingDirectory=/opt/vox
EnvironmentFile=/opt/vox/.env
ExecStart=/opt/vox/.venv/bin/python /opt/vox/maintenance.py
